# conf/httpbin.service
[Unit]
Description=chinaran/go-httpbin

[Service]
ExecStart=/bin/bash -c "HOSTNAME=$$(hostname) go-httpbin"

[Install]
WantedBy=multi-user.target
